# files/supervisor/kibana.ini
[program:kibana]
command=/opt/kibana/bin/kibana -c /opt/kibana/config/kibana.yml
autostart=true